# dots/.config/systemd/user/chromix-too.service
[Unit]
Description=Chromix too server

[Service]
ExecStart=/%h/.npm-global/bin/chromix-too-server

