# dots/.config/systemd/user/projects-sync.service
[Unit]
Description=Sync projects
After=network.target
Wants=network.target

[Service]
ExecStart=%h/.config/systemd/sync-projects.sh

