# dots/.config/systemd/user/twitch-local-commands.service
[Unit]
Description=Chromix too server
Requires=chromix-too.service

[Service]
Environment="PATH=/usr/bin:/home/mnbv/.npm-global/bin/"
EnvironmentFile=%h/.private/google
WorkingDirectory=%h/projects/nglgzz/twitch
ExecStart=/usr/bin/node server.js
ExecStop=/usr/bin/systemctl stop --user chromix-too.service

